# configfiles/local.d/reflector.start
#! /bin/bash
reflector --save /etc/pacman.d/mirrorlist-arch --country Germany --protocol https --latest 5